# systemd/backup-db.timer
[Unit]
Description=Run DB backup every 2 days at 02:00

[Timer]
OnCalendar=*-*-*/2 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
