# systemd/innopapp-api.service
[Unit]
Description=Innovation Project API managed by PM2
Documentation=man:startServer.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xavier
WorkingDirectory=/home/xavier/COS30049InnoProj
Environment=NODE_ENV=production
ExecStart=/usr/bin/env pm2-runtime start src/server.js --name innopapp-api --update-env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target